# web/bonus-calculator.service
[Unit]
Description=Bonus Calculator Web Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/bonus_calculator
ExecStart=/usr/bin/python3 /var/www/bonus_calculator/server.py 8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bonus-web

# 环境变量
Environment=PYTHONUNBUFFERED=1

# 资源限制
MemoryMax=256M
TasksMax=10

[Install]
WantedBy=multi-user.target